# backend/supervisord.conf
# =========================================
# Supervisord Configuration for NovaSight Backend
# =========================================
# Manages Flask API, Dagster Webserver, and Dagster Daemon

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# =========================================
# Flask API Server
# =========================================
[program:flask]
command=gunicorn --bind 0.0.0.0:5000 --workers 2 --threads 4 --timeout 120 "app:create_app()"
directory=/app
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/flask.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/flask-error.log
stderr_logfile_maxbytes=10MB
environment=FLASK_APP="app",FLASK_ENV="%(ENV_FLASK_ENV)s"
priority=100

# =========================================
# Dagster Webserver
# =========================================
[program:dagster-webserver]
command=dagster-webserver -h 0.0.0.0 -p 3000 -w /app/orchestration/workspace.yaml
directory=/app
autostart=true
autorestart=true
startsecs=15
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/dagster-webserver.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/dagster-webserver-error.log
stderr_logfile_maxbytes=10MB
environment=DAGSTER_HOME="/app/orchestration"
priority=200

# =========================================
# Dagster Daemon (Scheduler, Sensors, Run Queue)
# =========================================
[program:dagster-daemon]
command=dagster-daemon run -w /app/orchestration/workspace.yaml
directory=/app
autostart=true
autorestart=true
startsecs=20
stopwaitsecs=60
stdout_logfile=/var/log/supervisor/dagster-daemon.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/dagster-daemon-error.log
stderr_logfile_maxbytes=10MB
environment=DAGSTER_HOME="/app/orchestration"
priority=300

# =========================================
# Group for all services
# =========================================
[group:novasight]
programs=flask,dagster-webserver,dagster-daemon
priority=999
